# supervisor_awslogs.conf
[program:awslogs]
command=/var/awslogs/bin/awslogs-agent-launcher.sh
autorestart=true
autostart=true
redirect_stderr=true
